# ops/com.visualradar.daily.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.visualradar.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "/绝对路径/Visual Radar" &amp;&amp; /usr/local/bin/pnpm daily</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/visual-radar-daily.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/visual-radar-daily-error.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
